# systemd/wifi-fix.service
[Unit]
Description=Wifi Fix

[Service]
Type=simple
ExecStart=/opt/wifi-fix.sh

[Install]
WantedBy=network.target
